# root/etc/s6-overlay/s6-rc.d/init-mod-st-mod-cuckoo-postinstall/up
/etc/s6-overlay/s6-rc.d/init-mod-st-mod-cuckoo/run
