# root/etc/s6-overlay/s6-rc.d/init-mod-st-mod-template-postinstall/run
#!/usr/bin/with-contenv bash

# Logic here
