# root/etc/s6-overlay/s6-rc.d/init-mod-st-mod-cuckoo-install/run
#!/usr/bin/with-contenv bash

## Ubuntu
if [ -f /usr/bin/apt ]; then
    echo "\
        wget \
        xz-utils" >> /mod-repo-packages-to-install.list

fi
