# root/etc/s6-overlay/s6-rc.d/svc-mod-cuckoo/run
#!/usr/bin/with-contenv bash

/bin/warp-svc
